# DEBIAN/postinst
#!/bin/bash

# sudo adduser --system counter
# sudo systemctl daemon-reload
# sudo systemctl enable counter.service
# sudo systemctl start counter.service
