# DEBIAN/postrm
#!/bin/bash

# sudo systemctl disable counter
# sudo systemctl daemon-reload

# if id "counter" &>/dev/null; then
#     sudo userdel counter
# fi
