# DEBIAN/prerm
#!/bin/bash

# sudo systemctl stop counter.service

# sudo systemctl disable counter.service
